# conf/systemd/user/mrw-lightning-xweather-fetch.service
# Normally DISABLED — wx-core owns Xweather (LaunchAgent). Enable only as emergency fallback.
[Unit]
Description=MRW Xweather lightning fetch (fallback on wx-i9; primary is wx-core)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/scott/wx/radar-foundry
ExecStart=/home/scott/wx/radar-foundry/bin/run_lightning_xweather_fetch_wx_i9.sh
Restart=always
RestartSec=5
StandardOutput=append:/tmp/lightning_xweather_fetch_wxi9.log
StandardError=append:/tmp/lightning_xweather_fetch_wxi9.log

[Install]
WantedBy=default.target
