# conf/launchd/com.mrw.radar_coordinator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mrw.radar_coordinator</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/scott/wx/radar-foundry/bin/run_radar_coordinator.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>/tmp/radar_coordinator_launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/radar_coordinator_launchd.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/scott/wx/radar-foundry</string>
</dict>
</plist>
