# conf/launchd/com.mrw.watchdog_radar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mrw.watchdog_radar</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/scott/wx/radar-foundry/bin/watchdog_radar_loops.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>StandardOutPath</key>
  <string>/tmp/mrw_watchdog_launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mrw_watchdog_launchd.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/scott/wx/radar-foundry</string>
</dict>
</plist>
