# conf/launchd/com.mrw.mrms_loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mrw.mrms_loop</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/scott/wx/radar-foundry/bin/run_mrms_loop.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>120</integer>
  <key>StandardOutPath</key>
  <string>/tmp/mrms_loop_launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mrms_loop_launchd.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/scott/wx/radar-foundry</string>
</dict>
</plist>
