# conf/launchd/com.mrw.lightning_xweather_fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mrw.lightning_xweather_fetch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/scott/wx/radar-foundry/bin/run_lightning_xweather_fetch.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/tmp/lightning_xweather_fetch.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/lightning_xweather_fetch.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/scott/wx/radar-foundry</string>
</dict>
</plist>
